# launchd/dev.pyrycode.pyry.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.pyrycode.pyry</string>

    <!-- Adjust these two paths for your install. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/pyry</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/YOU/pyry-workspace</string>

    <!-- Ensure `claude` is on PATH. Adjust as needed. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/pyry.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/pyry.err.log</string>
</dict>
</plist>
